# ash-arch/backup.timer
[Unit]
Description="Timer for automatic backup"

[Timer]
Unit=backup.service
OnBootSec=20m
OnUnitInactiveSec=1d

[Install]
WantedBy=timers.target
